# deploy/patent-client-agents-mcp.service
[Unit]
Description=patent-client-agents MCP Server
After=network.target

[Service]
Type=simple
User=patent-client-agents
Group=patent-client-agents
WorkingDirectory=/opt/patent-client-agents
EnvironmentFile=/opt/patent-client-agents/.env
ExecStart=/opt/patent-client-agents/.venv/bin/fastmcp run src/patent_client_agents/mcp/server.py:mcp --transport streamable-http --stateless --host 127.0.0.1 --port 8002
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=patent-client-agents-mcp

[Install]
WantedBy=multi-user.target
